# phi3.5_llamacpp_service/llama-phi35.service
[Unit]
Description=Phi-3.5 llama.cpp Service
After=network.target

[Service]
Type=simple
# Update these to your runtime account/group if needed.
User=xptwb
Group=grpxptwb
WorkingDirectory=/opt/phi35_llamacpp

# Optional env file generated by installer flow.
EnvironmentFile=-/opt/phi35_llamacpp/config/phi35.env

# Keep this explicit and simple for airgapped environments.
ExecStart=/usr/bin/env LD_LIBRARY_PATH=/opt/phi35_llamacpp/lib /opt/phi35_llamacpp/bin/llama-server --model /opt/phi35_llamacpp/models/Phi-3.5-mini-instruct-Q4_K_M.gguf --host 127.0.0.1 --port 8000 --threads 4 --threads-batch 4 --parallel 1 --ctx-size 16384 --n-predict 1024 --cache-type-k q8_0 --cache-type-v q8_0 --metrics --no-webui --cont-batching

Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed

StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-phi35

[Install]
WantedBy=multi-user.target
