# onprem_llamacpp_service/systemd/llamacpp.service
[Unit]
Description=llama.cpp Inference Server (PoC)
Documentation=https://github.com/ggml-org/llama.cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llamacpp
Group=llamacpp
WorkingDirectory=/opt/llamacpp
EnvironmentFile=-/etc/llamacpp/llamacpp.env
ExecStart=/usr/bin/bash -lc '/usr/local/bin/llama-server --model "${LLAMA_MODEL_PATH}" --host "${LLAMA_HOST}" --port "${LLAMA_PORT}" --threads "${LLAMA_THREADS}" --threads-batch "${LLAMA_THREADS_BATCH}" --parallel "${LLAMA_PARALLEL}" --ctx-size "${LLAMA_CTX_SIZE}" --n-predict "${LLAMA_DEFAULT_MAX_TOKENS}" --cache-type-k "${LLAMA_CACHE_TYPE_K}" --cache-type-v "${LLAMA_CACHE_TYPE_V}" ${LLAMA_CONT_BATCHING_FLAG} ${LLAMA_MMAP_FLAG} ${LLAMA_MLOCK_FLAG} --metrics --no-webui ${LLAMA_EXTRA_ARGS}'
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
ProtectHome=yes
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llamacpp

[Install]
WantedBy=multi-user.target
