# llm_notable_analysis_onprem/systemd/vllm.service
[Unit]
Description=vLLM Inference Server (gpt-oss-20b)
Documentation=https://docs.vllm.ai/
After=network.target

[Service]
Type=simple
User=vllm
Group=vllm

# Working directory for model cache
WorkingDirectory=/opt/vllm

# vLLM server command
# Adjust --model path to local model directory
# Adjust --gpu-memory-utilization as needed (0.9 = 90% of GPU memory)
#
# SECURITY NOTE: --trust-remote-code is DISABLED by default for hardening.
# This flag allows executing arbitrary Python code bundled with model artifacts.
# Only enable if your model requires custom architecture code AND you have
# verified model artifacts via checksum from a trusted offline import process.
# To enable: add --trust-remote-code to ExecStart below.
#
ExecStart=/opt/vllm/venv/bin/python -m vllm.entrypoints.openai.api_server \
    --model /opt/models/gpt-oss-20b \
    --served-model-name gpt-oss-20b \
    --host 127.0.0.1 \
    --port 8000 \
    --gpu-memory-utilization 0.9 \
    --max-model-len 131072 \
    --dtype auto

# Graceful shutdown
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Restart policy
Restart=on-failure
RestartSec=30

# Environment for CUDA
Environment="CUDA_VISIBLE_DEVICES=0"

# Security hardening
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077
ProtectHome=yes
PrivateTmp=yes

# Network hardening: vLLM is local-only and should not need off-host egress.
IPAddressDeny=any
IPAddressAllow=127.0.0.1

# Resource limits
LimitNOFILE=65535
LimitMEMLOCK=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm

[Install]
WantedBy=multi-user.target

