# glab_vllm_litellm_kb_host_setup/onprem_kb_indexer/systemd/kb-rebuild.timer
[Unit]
Description=Daily GLAB Knowledge Base Rebuild

[Timer]
OnBootSec=10m
OnUnitActiveSec=24h
RandomizedDelaySec=5m
Persistent=true
Unit=kb-rebuild.service

[Install]
WantedBy=timers.target
